# web/docker/supervisord.conf
[supervisord]
nodaemon=true

[program:Flask]
environment=prometheus_multiproc_dir=/deploy/metrics
directory=/deploy
command=/deploy/manage.py runserver --host 0.0.0.0 --port 5000

